# deploy/systemd/medical-pbl.service
[Unit]
Description=Medical PBL Node Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/medical-pbl
Environment=NODE_ENV=production
EnvironmentFile=/opt/medical-pbl/.env
ExecStart=/usr/bin/node /opt/medical-pbl/server.js
Restart=always
RestartSec=5
TimeoutStopSec=15
KillSignal=SIGINT

# 日志交给 journald
StandardOutput=journal
StandardError=journal

# 资源与安全（最小可用）
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535
ReadWritePaths=/opt/medical-pbl-data /opt/medical-pbl

[Install]
WantedBy=multi-user.target
